# roles/galaxyprojectdotorg.cloudman-image/files/cloudman.conf
description     "Start CloudMan contextualization script"

start on runlevel [2345]

task
exec python /usr/bin/ec2autorun.py 2> /usr/bin/ec2autorun.log
