# roles/cloudman_image/files/supervisor/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
umask=022
autostart=true
autorestart=true
startsecs=5
exitcodes=0
priority=510
